# setup/udev_rules/88-cyusb.rules
# Vendored from prophesee-ai/openeb (Apache-2.0):
#   hal_psee_plugins/resources/rules/88-cyusb.rules
# Grants non-root access to the Cypress FX3 USB bridge used by the EVK4
# (vendor 04b4). The cy_renumerate.sh clause only matters with Cypress' own
# SDK; udev logs a harmless warning about the missing script otherwise.
#
# Cypress USB driver for FX2 and FX3 (C) Cypress Semiconductor Corporation / ATR-LABS
# Rules written by V. Radhakrishnan ( rk@atr-labs.com )
# Cypress USB vendor ID = 0x04b4
KERNEL=="*", SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_device", ACTION=="add", ATTR{idVendor}=="04b4", MODE="666", TAG="cyusb_dev", RUN+="/usr/local/bin/cy_renumerate.sh A"
KERNEL=="*", SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_device", ACTION=="remove", TAG=="cyusb_dev", RUN+="/usr/local/bin/cy_renumerate.sh R"
